# rel/files/installation/init.d/akkoma
#!/sbin/openrc-run

supervisor=supervise-daemon

# Requires OpenRC >= 0.35
directory=/opt/akkoma

command=/opt/akkoma/bin/pleroma
command_args="start"
command_user=akkoma
command_background=1

# Ask process to terminate within 30 seconds, otherwise kill it
retry="SIGTERM/30/SIGKILL/5"

pidfile="/var/run/akkoma.pid"

depend() {
    want nginx
    need postgresql
}
